# systemd/hermes-knowledge-sync.service
[Unit]
Description=Hermes Knowledge Auto-Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/klausd/Developer/scripts/hermes-knowledge-sync.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
